# deploy/photomind-frontend.service
[Unit]
Description=PhotoMind Frontend (Next.js / Bun)
Documentation=https://github.com/karthiknitt/photomind
After=network.target photomind-daemon.service photomind-bridge.service

[Service]
Type=simple
User=karthik
WorkingDirectory=/home/karthik/projects/PhotoMind/frontend

# Production environment
Environment=NODE_ENV=production
Environment=PORT=3003
Environment=NEXT_TELEMETRY_DISABLED=1

# Secrets and paths (edit this file on the VPS — never commit it)
EnvironmentFile=/home/karthik/photomind/.env.production

ExecStart=/usr/local/bin/bun run start
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=3

# Output
StandardOutput=journal
StandardError=journal
SyslogIdentifier=photomind-frontend

[Install]
WantedBy=multi-user.target
